# packages/nadybot@.service
[Unit]
Description=Nadybot Anarchy Online chatbot
After=network.target

[Service]
Type=simple
User=nadybot
Group=nadybot
WorkingDirectory=/opt/nadybot
ExecStart=/usr/bin/env php \
	/opt/nadybot/main.php \
	"/etc/nadybot/%i.php"
Restart=always
RestartPreventExitStatus=10
RestartSec=5
SyslogIdentifier=nadybot-%i

[Install]
WantedBy=multi-user.target